# deb/src/deb/control/postinst
#!/bin/sh

set -e

update-rc.d ${build.finalName} defaults
service ${build.finalName} start

exit 0