# deb/src/deb/control/postrm
#!/bin/sh

set -e

update-rc.d -f ${build.finalName} remove

if [ -e /var/log/${build.finalName}.log ]; then
  rm /var/log/${build.finalName}.log
fi

exit 0